# containers/xfce4/rootfs/home/user/.config/sv/xfce4/run
#!/usr/bin/env /lib/runit/invoke-run

exec 2>&1

if [ -e /etc/runit/verbose ]; then
    echo "invoke-run: starting ${PWD##*/}"
fi

SVDIR=$HOME/.local/service sv start xwayland || exit 1

export WAYLAND_DISPLAY=
export DISPLAY=:0

cd $HOME

exec /usr/bin/dbus-run-session /usr/bin/xfce4-session
