# containers/xfce4/rootfs/home/user/.config/sv/xwayland/run
#!/usr/bin/env /lib/runit/invoke-run

exec 2>&1

if [ -e /etc/runit/verbose ]; then
    echo "invoke-run: starting ${PWD##*/}"
fi

exec /usr/bin/Xwayland :0
