# containers/weston/rootfs/home/user/.config/sv/weston/run
#!/usr/bin/env /lib/runit/invoke-run

exec 2>&1

if [ -e /etc/runit/verbose ]; then
    echo "invoke-run: starting ${PWD##*/}"
fi

export XDG_SESSION_ID=2

exec /usr/bin/weston -Swayland-2 --width 1280 --height 720
